# vidality_resume_service/vidality_resume_service.ini
[program:vidality_resume_service]
user=marc
command=/home/marc/bin/python /home/marc/python_test_apps/vidality_resume_service/filter_service.py
directory=/home/marc/python_test_apps/vidality_resume_service
autostart=true
autorestart=true
logfile=log.txt
