# vidality_linkedin_service/vidality_linkedin_service.ini
[program:vidality_linkedin_service]
user=marc
command=/home/marc/bin/python /home/marc/python_test_apps/vidality_linkedin_service/server.py
directory=/home/marc/python_test_apps/vidality_linkedin_service
autostart=true
autorestart=true
logfile=log.txt
